# numa_scripts/xen_setup/guest_setup/cpu-online.rules
SUBSYSTEM=="cpu", ACTION=="add", RUN+="/bin/sh -c 'echo 1 > %S%p/online'"
